# crontab
# run every 3 hours

* */3 * * * /appsrc/run_job.sh >> /var/log/cron.log 2>&1


